# act_take/holy_image_2015_09_09/imx5/system/dropbear.service
[Unit]
Description=ssh deamon

[Service]
Type=forking
ExecStart=/lib/systemd/systemd-dropbear start
ExecStop=/lib/systemd/systemd-dropbear stop

[Install]
WantedBy=multi-user.target

