# act_take/holy_image_2015_09_09/imx5/system/dvi-enable.service
[Unit]
Description=Enable DVI

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/lib/systemd/systemd-dvi-enable start
ExecStop=/lib/systemd/systemd-dvi-enable stop

[Install]
WantedBy=multi-user.target
